# src/vision/camera/systemd/vision_camera_rpi.service
[Unit]
Description=SmartMower Vision Camera (Raspberry Pi)
Documentation=https://github.com/yourorg/smartmower
After=network.target mosquitto.service
Requires=mosquitto.service

[Service]
Type=simple
User=smartmower
Group=video  # Add video group for camera access
WorkingDirectory=/opt/smartmower
ExecStart=/opt/smartmower/bin/vision_camera_rpi
Restart=always
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vision-camera-rpi

# Security hardening
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
MemoryDenyWriteExecute=true
LockPersonality=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
CapabilityBoundingSet=
DeviceAllow=/dev/video0 rw
DeviceAllow=/dev/vchiq rw
DeviceAllow=char-alsa rw
DeviceAllow=char-drm rw

# Resource limits
LimitNOFILE=65536
LimitNPROC=512
LimitMEMLOCK=64M

# Directory access
ReadWritePaths=/opt/smartmower/log /opt/smartmower/data/vision
ReadOnlyPaths=/
ReadOnlyPaths=-/etc
ReadOnlyPaths=-/usr
ReadOnlyPaths=-/lib
ReadOnlyPaths=-/lib64
ReadOnlyPaths=-/bin
ReadOnlyPaths=-/sbin
ReadWritePaths=-/tmp

# Environment
Environment=GST_DEBUG=2
Environment=GST_DEBUG_NO_COLOR=1

[Install]
WantedBy=multi-user.target
