# src/vision/camera/systemd/vision_camera.service
[Unit]
Description=Smart Mower Vision Camera Service
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=smartmower
Group=smartmower
WorkingDirectory=/opt/smartmower
ExecStart=/opt/smartmower/bin/vision_camera
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/opt/smartmower/log /opt/smartmower/data/vision

[Install]
WantedBy=multi-user.target
