# src/_migrato/pico/systemd/pico_bridge.service
[Unit]
Description=Smart Mower Pico MQTT Bridge
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=smartmower
Group=smartmower
WorkingDirectory=/opt/smartmower
ExecStart=/opt/smartmower/bin/pico_bridge /opt/smartmower/etc/config/robot_config.json
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=MQTT_HOST=localhost
Environment=MQTT_PORT=1883

[Install]
WantedBy=multi-user.target
