# src/gps/systemd/gps_bridge.service
[Unit]
Description=GPS Bridge for Smart Mower
After=network.target mosquitto.service

[Service]
Type=simple
User=smartmower
Group=smartmower
WorkingDirectory=/opt/smartmower/bin
ExecStart=/opt/smartmower/bin/gps_bridge /opt/smartmower/etc/config/robot_config.json
Restart=always
RestartSec=10

# Security options
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gps-bridge

[Install]
WantedBy=multi-user.target
