# src/bridge/gps/systemd/gps_bridge.service
[Unit]
Description=SmartMower GPS Bridge Service
After=network.target mosquitto.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/smartmower/bin
ExecStart=/opt/smartmower/bin/gps_bridge
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=gps_bridge
Environment="LD_LIBRARY_PATH=/usr/local/lib"

[Install]
WantedBy=multi-user.target
