# src/_damigrare/state_machine/systemd/state_machine.service
[Unit]
Description=Smart Mower State Machine Controller
After=network.target mosquitto.service fusion_sensor.service
Wants=mosquitto.service fusion_sensor.service
Requires=mosquitto.service

[Service]
Type=simple
User=smartmower
Group=smartmower
WorkingDirectory=/opt/smartmower
ExecStart=/opt/smartmower/bin/state_machine
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=MQTT_HOST=localhost
Environment=MQTT_PORT=1883

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/smartmower/log /var/lib/smartmower

[Install]
WantedBy=multi-user.target
