# src/safety/smartmower-safety.service
[Unit]
Description=SmartMower Safety Supervisor
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=smartmower
Group=smartmower
ExecStart=/opt/smartmower/bin/safety_supervisor
Restart=always
RestartSec=2
WorkingDirectory=/opt/smartmower

[Install]
WantedBy=multi-user.target
