# src/_damigrare/config/systemd/parameters.service
[Unit]
Description=Smart Mower Parameter Management Service
After=network.target mosquitto.service
Wants=network.target
Requires=mosquitto.service

[Service]
Type=simple
User=smartmower
Group=smartmower
WorkingDirectory=/opt/smartmower/bin
ExecStart=/usr/bin/python3 /opt/smartmower/bin/parameter_service.py --config /opt/smartmower/etc/config/robot_config.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smartmower-parameters

# Environment variables
Environment=PYTHONPATH=/opt/smartmower/lib
Environment=MQTT_BROKER=localhost
Environment=MQTT_PORT=1883

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/smartmower/etc/config /opt/smartmower/log

[Install]
WantedBy=multi-user.target
